# setup/files/fieldrig-server.service
[Unit]
Description=FieldRig vehicle terminal server
After=pipewire.service

[Service]
ExecStart=%h/.local/share/fieldrig/venv/bin/python -m fieldrig --host 127.0.0.1 --port 8000
Environment=FIELDRIG_ENABLE_SHUTDOWN=1
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
